# 24_backup/borg/borgbackup.timer
[Unit]
Description=Borg backup timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=borgbackup.service

[Install]
WantedBy=multi-user.target
